# cisco-secure-client-linux64-predeploy-k9/cisco-secure-client-linux64-5.1.11.388/posture/ciscod.service
[Unit] 
Description=Cisco Secure Client - Posture Daemon 

[Service] 
Type=forking
Restart=on-failure
ExecStart=/opt/cisco/secureclient/securefirewallposture/bin/ciscod
PIDFile=/var/run/ciscod.pid
KillMode=process
EnvironmentFile=/etc/environment

[Install] 
WantedBy=multi-user.target